# r2_control.service
[Unit]
Description=R2 Control
After=network.target
 
[Service]
Type=simple
User=root
Group=root
ExecStart=/opt/r2_control/venv/bin/gunicorn -w 1 -b 0.0.0.0:5000 main:app
WorkingDirectory=/opt/r2_control
StandardOutput=syslog
StandardError=syslog
Restart=on-failure
 
[Install]
WantedBy=network.target
Alias=r2_control.service


