# r2_control.service
[Unit]
Description=R2 Control
 
[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/python -O -u /home/pi/r2_control/main.py
WorkingDirectory=/home/pi/r2_control
StandardOutput=syslog
StandardError=syslog
Restart=on-failure
 
[Install]
WantedBy=multi-user.target
Alias=r2_control.service


